# dist/bartoc/bartoc-logrotate.timer
[Unit]
Description=Daily bartoc log rotation

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
