# dist/bartoc/bartoc-logrotate.service
[Unit]
Description=Rotate bartoc logs

[Service]
Type=oneshot
ExecStart=/usr/lib/bartoc/bartoc-logrotate
